# vlm/systemd/vlm-daemon.service
[Unit]
Description=LeKiwi VLM daemon — camera + caption HTTP API (three-state power model)
After=llama-server.service
Wants=llama-server.service

[Service]
Type=simple
WorkingDirectory=%h/work/lekiwi-jatson-orin/vlm
# Runs even when llama-server is down (captions degrade gracefully; /frame.jpg
# stays up since the camera is independent of the VLM backend).
ExecStart=%h/work/lekiwi-jatson-orin/vlm/.venv/bin/python \
    %h/work/lekiwi-jatson-orin/vlm/daemon.py
Restart=on-failure
RestartSec=3
# Make sure ~/.local/bin/ffmpeg is resolvable.
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
