# board/home/jatson/.config/systemd/user/hermes-gateway-robot.service
[Unit]
Description=Hermes Agent Gateway - Messaging Platform Integration
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/home/jatson/.hermes/hermes-agent/venv/bin/python -m hermes_cli.main --profile robot gateway run
WorkingDirectory=/home/jatson/.hermes/profiles/robot
Environment="PATH=/home/jatson/.hermes/hermes-agent/venv/bin:/home/jatson/.hermes/hermes-agent/node_modules/.bin:/home/jatson/.nvm/versions/node/v22.23.1/bin:/home/jatson/.local/bin:/home/jatson/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="VIRTUAL_ENV=/home/jatson/.hermes/hermes-agent/venv"
Environment="HERMES_HOME=/home/jatson/.hermes/profiles/robot"
Restart=always
RestartSec=5
RestartForceExitStatus=75
RestartPreventExitStatus=78
KillMode=mixed
KillSignal=SIGTERM
ExecReload=/bin/kill -USR1 $MAINPID
ExecStopPost=-/home/jatson/.hermes/hermes-agent/venv/bin/python -m gateway.cgroup_cleanup
TimeoutStopSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
