# board/etc/systemd/system/pad_teleop.service
# LeKiwi gamepad teleop daemon — evdev pad -> ZMQ tcp://127.0.0.1:5555.
# Independent of the GUI; hotplug-safe (waits for the pad if absent).
[Unit]
Description=LeKiwi gamepad teleop (evdev -> ZMQ)
After=base_host.service
Wants=base_host.service
# Never stop retrying — a Bluetooth pad may pair minutes after boot, and we
# were once found dead (inactive) with Restart alone. 0 = no start-rate limit.
StartLimitIntervalSec=0

[Service]
User=jatson
ExecStart=/home/jatson/miniconda3/envs/lerobot/bin/python /home/jatson/pad_teleop.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
